# crontab/crontabrc.symlink
MAILTO=""
#@weekly
0 0 * * 0 . $HOME/.bash_profile; $HOME/.dotfiles/crontab/dotToGit.sh

